# buildroot/package/snapdog-client/snapdog-client.service
[Unit]
Description=SnapDog Client (Snapcast audio receiver)
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/snapdog-client
ExecStart=/usr/bin/snapdog-client $SNAPDOG_CLIENT_ARGS
Restart=on-failure
RestartSec=5
DynamicUser=yes
SupplementaryGroups=audio
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
