# buildroot/package/snapdog-updater/boot-guard.service
[Unit]
Description=Boot Guard (OTA rollback protection)
DefaultDependencies=no
Before=snapdog-ctrl.service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/opt/snapdog/bin/boot-guard
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
